# bin/spc.service
# https://www.freedesktop.org/software/systemd/man/systemd.service.html
[Unit]
Description=Pironman U1 Auto Control
After=multi-user.target # Need to start last to avoid gpio being occupied

[Service]
# ExecStartPre=/bin/sleep 15
Type=forking
WorkingDirectory=/opt/spc
ExecStart=python3 /opt/spc/spc_server start
PrivateTmp=True

[Install]
WantedBy=multi-user.target
